# crontab
20 10 * * * /usr/bin/python3 /Users/yourname/health_sync/dingtalk_report.py